# meta-hdsc/meta-hs9216/recipes-core/systemd/obmc-targets/obmc-chassis-reset@.target
[Unit]
Description=HS9216 Chassis%i Warm Reset
Wants=obmc-standby.target
After=obmc-standby.target
Conflicts=obmc-chassis-poweroff@%i.target
RefuseManualStop=no
