# meta-hdsc/meta-hs2500/recipes-phosphor/gpio/fan-test/fan-test-button-pressed.service
[Unit]
Description=Fan Test Button Pressed Service

[Service]
Type=oneshot
ExecStart=/usr/bin/env toggle_fan_speed.sh
SyslogIdentifier=fan-test-button
